# ansible/roles/reunhangout/templates/reunhangout-workers.target
# This is a target for controlling launching multiple reunhangout workers at once.

[Install]
WantedBy=multi-user.target
